# crates/agent-core/debian/sentinel-agent.service
[Unit]
Description=Sentinel GRC Agent - Agent de conformite endpoint
Documentation=https://docs.sentinel-grc.com
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/sentinel-grc/bin/sentinel-agent --service
Restart=always
RestartSec=10
WatchdogSec=120

# Dedicated service user
User=sentinel-grc
Group=sentinel-grc
WorkingDirectory=/opt/sentinel-grc
Environment=RUST_LOG=info

# Filesystem access
ReadWritePaths=/var/lib/sentinel-grc /var/log/sentinel-grc /etc/sentinel

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
RemoveIPC=true
SystemCallArchitectures=native

# Capabilities (network discovery, filesystem scanning)
CapabilityBoundingSet=CAP_NET_RAW CAP_DAC_READ_SEARCH
AmbientCapabilities=CAP_NET_RAW CAP_DAC_READ_SEARCH

# System call filter (allowlist)
SystemCallFilter=@system-service
SystemCallFilter=~@mount @reboot @swap @debug @obsolete

# Resource limits
MemoryMax=150M
MemoryHigh=100M
CPUQuota=10%
TasksMax=64
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
